# infrastructure/ansible/roles/restic/templates/plex-backup.timer.j2
[Unit]
Description=Daily Plex Backup Timer

[Timer]
OnCalendar={{ backup_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
